# home/.config/autostart/dropboxstart.desktop
[Desktop Entry]
Type=Application
Name=dropboxstart
Exec=/home/alex/bin/db
