# home/.config/autostart/updatecheck.desktop
[Desktop Entry]
Type=Application
Name=updatecheck
Exec=/usr/bin/aptdcon --refresh
